# miniflux/podman/db.container
[Unit]
Description=Miniflux PostgreSQL Database
Wants=miniflux.pod
After=miniflux.pod

[Container]
# Use the PostgreSQL 17 image
Image=postgres:17

# Set container name
ContainerName=miniflux-db

# Environment variables from environment file
EnvironmentFile=/etc/miniflux/miniflux.env
Environment=POSTGRES_USER=${PG_USER}
Environment=POSTGRES_PASSWORD=${PG_PASS}
Environment=POSTGRES_DB=${PG_DB}

# Health check
HealthCmd=/usr/bin/pg_isready -U ${PG_USER}
HealthInterval=10s
HealthStartupPeriod=30s

# Use named volume
Volume=miniflux_db:/var/lib/postgresql/data

# Add to the pod
Pod=miniflux-pod

[Install]
WantedBy=miniflux.pod
